# supervisor.dev.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
loglevel=debug

[program:frontend-dev]
directory=/app/ui
command=bash -c "cd /app/ui && WATCHPACK_POLLING=true pnpm dev"
autostart=true
autorestart=true
stderr_logfile=/var/log/frontend.err.log
stdout_logfile=/var/log/frontend.out.log
environment=NODE_ENV="development",PORT="%(ENV_FRONTEND_PORT)s",LOG_LEVEL="debug",NEXT_TELEMETRY_DISABLED="1",WATCHPACK_POLLING="true"
stopasgroup=true
killasgroup=true
startsecs=10
startretries=3

[program:backend-dev]
directory=/app
command=uvicorn backend.app:app --host 0.0.0.0 --port %(ENV_BACKEND_PORT)s --reload --log-level debug --reload-dir /app/backend --workers 1
autostart=true
autorestart=true
stderr_logfile=/var/log/backend.err.log
stdout_logfile=/var/log/backend.out.log
stopasgroup=true
killasgroup=true 